# autostart/redshift.desktop
[Desktop Entry]
Type=Application
Name=Redshift
Terminal=false
Exec=bash /home/memo/redshift.sh
NoDisplay=true

