# autostart/redshift.desktop
[Desktop Entry]
Type=Application
Name=Redshift
Terminal=false
Exec=bash ~/redshift.sh
NoDisplay=true

